# 02_Projects/task_bridge/com.woosdom.hands3-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.woosdom.hands3-runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>-u</string>
    <string>/Users/woosung/Desktop/Dev/Woosdom_Brain/02_Projects/task_bridge/hands3_runner.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/Users/woosung/.npm-global/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>SSL_CERT_FILE</key>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/lib/python3.14/site-packages/certifi/cacert.pem</string>
    <key>REDIS_HOST</key>
    <string>127.0.0.1</string>
    <key>REDIS_PORT</key>
    <string>6379</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/woosung/Desktop/Dev/Woosdom_Brain/02_Projects/task_bridge</string>
  <key>StandardOutPath</key>
  <string>/tmp/woosdom-hands3-runner.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/woosdom-hands3-runner-err.log</string>
</dict>
</plist>
